# template/services/node.service
[Unit]
Description=Service For WebApp
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
WorkingDirectory=/var/www/athios.noahvarghese.me
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
